# qr3_78.conf
[program:qr3_78]
command=python3 /home/curator3/qr3_78/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_78/test.err.log
stdout_logfile=/home/curator3/qr3_78/test.err.log
